# CortanaDisplay.service
[Unit]
Description=Cortana Display
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/CortanaDisplay/CortanaDisplay.py
Restart=on-abort

[Install]
WantedBy=multi-user.target
